# scripts/launchd/com.lpkralv.ccmanager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lpkralv.ccmanager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Volumes/Sheridan/sandbox/ClaudeCodeManager/dist/server/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/Sheridan/sandbox/ClaudeCodeManager</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3000</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/Sheridan/sandbox/ClaudeCodeManager/logs/ccmanager.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/Sheridan/sandbox/ClaudeCodeManager/logs/ccmanager-error.log</string>
</dict>
</plist>
